# soilguard.service
[Unit]
Description=SoilGuard Portal - Soil Quality & Agricultural Monitor
Documentation=https://github.com/fivepanelhat/SoilGuard-Portal
After=network.target ollama.service mosquitto.service
Requires=ollama.service

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/SoilGuard-Portal
Environment="PATH=/home/pi/SoilGuard-Portal/venv/bin"
Environment="PYTHONUNBUFFERED=1"

# Main executable
ExecStart=/home/pi/SoilGuard-Portal/venv/bin/python main.py

# Auto-restart on failure
Restart=on-failure
RestartSec=10

# Output handling (send to systemd journal)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=soilguard

# Resource limits (prevent runaway processes on RPi 5)
MemoryMax=4G
CPUQuota=80%
MemoryAccounting=yes

# Keep environment clean
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
